# files/sensu-gem/upstart/sensu-client.conf
# Upstart: /etc/init/sensu-client.conf

description "sensu client"

start on (remote-filesystems and net-device-up)
stop on [!12345]

respawn

env PROG=sensu-service
env USER=sensu

chdir /opt/sensu/bin

pre-start script
  test -x ./$PROG || { stop; exit 0; }
end script

setuid $USER
exec ./$PROG start client

