# config/templates/sensu-gem/rcd/sensu-service.erb
#!/bin/sh

# PROVIDE: <%= service_name %>
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="<%= service_name %>"
rcvar="<%= service_name %>_enable"

command="/opt/sensu/bin/sensu-service"
command_args="start <%= service_shortname %> fork"
command_interpreter="/opt/sensu/embedded/bin/ruby"

pidfile="/var/run/sensu/sensu-<%= service_shortname %>.pid"
procname="/opt/sensu/bin/sensu-<%= service_shortname %>"

extra_commands="validate_config"

start_precmd="start_precmd"
restart_precmd="restart_precmd"
validate_config_cmd="<%= service_name %>_validate_config"

<%= service_name %>_chroot="/"
<%= service_name %>_user="sensu"

<%= service_name %>_validate_config() {
  echo "Validating sensu-<%= service_shortname %> configuration"
  su -m $sensu_client_user -c "${command} validate <%= service_shortname %>"
  if [ $? != 0 ]; then
    return 1
  fi
}

start_precmd() {
  sensu_client_validate_config
}

restart_precmd() {
  sensu_client_validate_config
}

load_rc_config $name
run_rc_command "$1"
